# deploy/dispatch.service
[Unit]
Description=Dispatch — kanban agent OS
After=network.target
ConditionPathExists=%h/dispatch/server.mjs

[Service]
# Edit WorkingDirectory and EnvironmentFile for your machine.
WorkingDirectory=%h/dispatch
ExecStart=/usr/bin/env node %h/dispatch/server.mjs
Restart=on-failure
RestartSec=5
KillMode=process
EnvironmentFile=-%h/.config/dispatch/dispatch.env
Environment=NODE_ENV=production
Environment=DISPATCH_PORT=4400
Environment=DISPATCH_DATA=%h/dispatch-data
# Optional for machine-specific URLs:
# Environment=DISPATCH_PUBLIC_URL=https://dispatch.example.com
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
